# roles/fileserver/templates/zfs_snapshot_metrics.service.j2
# Managed by Ansible - do not edit manually
[Unit]
Description=ZFS Snapshot Metrics Collector for Prometheus Textfile Collector
After=zfs.target
Wants=zfs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zfs_snapshot_metrics.sh
User=root
