# roles/fileserver/templates/zfs_snapshot_metrics.timer.j2
# Managed by Ansible - do not edit manually
[Unit]
Description=ZFS Snapshot Metrics Collection Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=zfs-snapshot-metrics.service

[Install]
WantedBy=timers.target
