# roles/base/templates/zfs_metrics.timer.j2
# Managed by Ansible - do not edit manually
[Unit]
Description=ZFS Metrics Collection Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=zfs-metrics.service

[Install]
WantedBy=timers.target
